# supervisord.conf
[supervisord]
nodaemon=true

[program:gnrdaemon]
command="/usr/local/bin/gnrdaemon"

[program:gnrwsgisite]
command="/usr/local/bin/gnrwsgi"